# deploy/systemd/captcha-solver.service.example
[Unit]
Description=GPTGrok2API local captcha solver
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__APP_DIR__/captcha-solver
Environment=HOME=__RUN_HOME__
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=TURNSTILE_HEADLESS=0
ExecStart=/usr/bin/xvfb-run -a --server-args="-screen 0 1920x1080x24" __APP_DIR__/captcha-solver/.venv/bin/python -m uvicorn server:app --host 127.0.0.1 --port 8877
Restart=on-failure
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGINT
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
