# deploy/systemd/chatgpt2api.service.example
[Unit]
Description=GPTGrok2API local service
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__APP_DIR__
Environment=HOME=__RUN_HOME__
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-__APP_DIR__/.env
ExecStart=__APP_DIR__/.venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGINT
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
