# deploy/launchd/com.chatgpt2api.captcha-solver.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chatgpt2api.captcha-solver</string>

  <key>ProgramArguments</key>
  <array>
    <string>__APP_DIR__/captcha-solver/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>server:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8877</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__APP_DIR__/captcha-solver</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SOLVER_ALLOW_PRIVATE</key>
    <string>1</string>
    <key>TURNSTILE_HEADLESS</key>
    <string>0</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__APP_DIR__/logs/captcha-solver.log</string>
  <key>StandardErrorPath</key>
  <string>__APP_DIR__/logs/captcha-solver.error.log</string>
</dict>
</plist>
